# supervisord.conf
[program:queue_listen]
command=php /var/www/gafers/artisan queue:listen
autostart=true
autorestart=true
[supervisord]