# supervisord.conf
[program:foo]
command=php artisan queue:listen

[supervisord]